# systemd/status-page.service
[Unit]
Description=Generate status page

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/noah/status-page/generate_status.py
