# systemd/status-page.timer
[Unit]
Description=Regenerate status page every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
